# bin/satellite.boot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
	"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>satellite.boot</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/student/dev/satellite/bin/start_satellite</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
